# backend/deploy/systemd/vrtiq-backup.timer
[Unit]
Description=Run vrtIQ database backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
Unit=vrtiq-backup.service

[Install]
WantedBy=timers.target
