# backend/deploy/systemd/vrtiq-prune.service
[Unit]
Description=vrtIQ storage pruning job

[Service]
Type=oneshot
User=vrtiq
Group=vrtiq
WorkingDirectory=/srv/vrtiq/app/backend
EnvironmentFile=/etc/vrtiq/backend.env
ExecStart=/usr/bin/bash /srv/vrtiq/app/backend/deploy/prune-storage.sh
