# backend/deploy/systemd/vrtiq-prune.timer
[Unit]
Description=Run vrtIQ storage pruning every 6 hours

[Timer]
OnCalendar=*-*-* *:00/6:00
Persistent=true
Unit=vrtiq-prune.service

[Install]
WantedBy=timers.target
